# ops/systemd/cmo-lens-facebook-sync.service
[Unit]
Description=CMO Lens Facebook channel metrics sync
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/ju/cmo-engine-openclaw
EnvironmentFile=-/home/ju/.config/cmo-engine-openclaw/lens-sync.env
ExecStart=/usr/bin/node scripts/cmo-lens-facebook-sync.mjs
