# ops/systemd/cmo-facebook-channel-sync.service
[Unit]
Description=CMO Product native Facebook channel metrics sync
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/home/ju/.config/cmo-engine-openclaw/dashboard.env
Environment=CMO_PUBLIC_APP_URL=http://127.0.0.1:3002
WorkingDirectory=/home/ju/cmo-engine-openclaw
ExecStart=/usr/bin/curl --fail --show-error --silent --max-time 300 \
  -X POST "${CMO_PUBLIC_APP_URL}/api/internal/lens/apps/holdstation-mini-app/social/facebook/sync" \
  -H "Authorization: Bearer ${CMO_LENS_INTERNAL_API_KEY}" \
  -H "Content-Type: application/json" \
  --data '{"rangeKey":"this_week","mode":"refresh_if_stale","trigger":"systemd_timer","dryRun":false}'

[Install]
WantedBy=multi-user.target
